# rootfs/usr/lib/systemd/system/winesapos-plasmalogin-health-check.service
[Unit]
Description=Check the status of PLM and invoke a recovery console if it fails

[Service]
Type=oneshot
TimeoutSec=5m
ExecStart=/usr/local/bin/winesapos-plasmalogin-health-check.sh

[Install]
WantedBy=multi-user.target
