# files/kdeconnect-indicator.service
[Unit]
Description=KDEConnect Indicator

[Service]
ExecStart=/usr/bin/kdeconnect-indicator

[Install]
WantedBy=graphical-session.target
